# deploy/systemd/sterlingdex-public.service
[Unit]
Description=SterlingDEX public API surface
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/sterling/sterlingdex-adapter
EnvironmentFile=/opt/sterling/sterlingdex-adapter/deploy/public.env
ExecStartPre=/usr/bin/env npm run build:surface
ExecStart=/usr/bin/env node server.mjs
Restart=always
RestartSec=3
User=sterling
Group=sterling

[Install]
WantedBy=multi-user.target
